# codesandbox-rust/scripts/codesandbox-rust.service
[Unit]
Description=Codesandbox Rust - Secure Code Execution API
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/path/to/codesandbox-rust
ExecStart=/path/to/codesandbox-rust/target/release/api
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codesandbox-rust

# Environment variables
Environment=ROOTFS_PATH=/path/to/codesandbox-rust/rootfs/rootfs.ext4
Environment=KERNEL_PATH=/path/to/codesandbox-rust/kernel/vmlinux
Environment=API_PORT=8000
Environment=RUST_LOG=info

# Security settings
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target